# deployment/systemd/twf-v3-tile-server.service
[Unit]
Description=TWF V3 Tile Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/twf_v3
Environment=PATH=/opt/twf_v3/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/opt/twf_v3/.venv/bin/uvicorn backend.app.services.tile_server:app --host 127.0.0.1 --port 8201 --workers 3
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
